# cinder/cinder.conf
[program:cinder-api]
command=/usr/bin/cinder-api --config-file=/etc/cinder/cinder.conf --log-file=/var/log/cinder/cinder-api.log
user=cinder

[program:cinder-scheduler]
command=/usr/bin/cinder-scheduler --config-file=/etc/cinder/cinder.conf --log-file=/var/log/cinder/cinder-scheduler.log
user=cinder

[program:cinder-backup]
command=/usr/bin/cinder-backup --config-file=/etc/cinder/cinder.conf --log-file=/var/log/cinder/cinder-backup.log
user=cinder

[program:cinder-volume]
command=/usr/bin/cinder-volume --config-file=/etc/cinder/cinder.conf --log-file=/var/log/cinder/cinder-volume.log
user=cinder
