# heat/heat.conf
[program:heat-api]
command=/usr/bin/heat-api
user=heat

[program:heat-api-cfn]
command=/usr/bin/heat-api-cfn
user=heat

[program:heat-engine]
command=/usr/bin/heat-engine
user=heat
